# docker-legacybaseos/filesystem/etc/systemd/system/firstboot.service
[Unit]
Description=Finish the installation.

[Service]
ExecStart=/bin/bash /usr/sbin/firstboot.sh

[Install]
WantedBy=multi-user.target

